# resources/ubuntu_net.service
[Unit]
After=sshd.service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/net_setup.sh
